# files/system/usr/lib/systemd/user/dms.service.d/override.conf
# Drop-in override for the RPM-shipped dms.service (avengemedia/dms COPR).
# Baked into the image via the `files` module (files/system -> /).
#
# Force the Qt platform theme to qt6ct for the DMS user session so DMS's Qt
# theming coexists with vanilla KDE settings, per the DMS wiki. Set both the
# generic and Qt6-specific vars; QT_QPA_PLATFORMTHEME_QT6 wins on Qt6 when present.
[Service]
Environment=QT_QPA_PLATFORMTHEME=qt6ct
Environment=QT_QPA_PLATFORMTHEME_QT6=qt6ct
