# files/system/etc/containers/systemd/users/litellm.container
# LiteLLM Proxy — unified gateway for multiple LLM providers
#
# This rootless user service runs the LiteLLM proxy (https://docs.litellm.ai) as a
# localhost-only API gateway, letting you manage OpenAI, Anthropic, Azure, and other
# provider APIs from a single endpoint. The proxy listens on http://127.0.0.1:4000
# and includes a web UI for managing keys and models interactively.
#
# SETUP (one-time, per machine):
#
# 1. Create ~/.config/litellm/config.yaml with your API keys and model list:
#
#    model_list:
#      - model_name: gpt-4
#        litellm_params:
#          model: openai/gpt-4
#          api_key: sk-proj-...
#      - model_name: claude
#        litellm_params:
#          model: anthropic/claude-opus-4-8
#          api_key: sk-ant-...
#
#    general_settings:
#      master_key: sk-1234-your-secret-master-key
#
#    The master_key is the bearer token clients use to access the proxy.
#    Keys stay local in your home directory, never in git or the image.
#
# 2. Start the service:
#
#    systemctl --user start litellm.service
#    systemctl --user enable litellm.service  # (optional: auto-start at login)
#
# 3. Test:
#
#    curl http://127.0.0.1:4000/health
#    # or open the web UI: http://127.0.0.1:4000
#
# The service will fail to start if ~/.config/litellm/config.yaml is missing or
# malformed — check logs with: journalctl --user -u litellm.service

[Unit]
Description=LiteLLM proxy gateway
Documentation=https://docs.litellm.ai
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/berriai/litellm:v1.93.0
ContainerName=litellm
PublishPort=127.0.0.1:4000:4000
Volume=%h/.config/litellm/config.yaml:/app/config.yaml:ro,Z
Exec=--config /app/config.yaml

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
