# files/system/etc/containers/systemd/users/metapi.container
# metapi — unified LLM gateway aggregating multiple relay/reseller panels
#
# This rootless user service runs metapi (https://github.com/cita-777/metapi),
# a "relay-of-relays" gateway that unifies multiple upstream LLM providers —
# reseller/aggregator panels (New API, One API, OneHub, Sub2API, etc.), direct
# provider APIs (OpenAI/Anthropic/Gemini), and OAuth logins (Codex/Claude/Gemini
# CLI) — behind a single endpoint and key. It exposes both OpenAI-compatible
# (/v1/chat/completions) and Anthropic-compatible (/v1/messages) surfaces, so
# claude-code and codex-cli can both point at it.
#
# The UI/gateway listens on http://127.0.0.1:4000 (localhost-only).
#
# SECRETS (one-time, per machine — never bake these into the image):
#
#   The three secrets below MUST be provided at runtime via a local systemd
#   drop-in (not committed to git). Generate strong random values:
#
#     mkdir -p ~/.config/systemd/user/metapi.service.d
#     cat > ~/.config/systemd/user/metapi.service.d/secrets.conf <<EOF
#     [Service]
#     Environment=AUTH_TOKEN=$(openssl rand -hex 32)
#     Environment=PROXY_TOKEN=sk-$(openssl rand -hex 32)
#     Environment=ACCOUNT_CREDENTIAL_SECRET=$(openssl rand -hex 32)
#     EOF
#     systemctl --user daemon-reload
#
#   - AUTH_TOKEN: admin login for the web UI (changeable in UI after first boot)
#   - PROXY_TOKEN: the key downstream clients (claude-code/codex-cli) authenticate
#     with when calling /v1/* (changeable in UI after first boot)
#   - ACCOUNT_CREDENTIAL_SECRET: encrypts stored upstream account credentials at rest
#
# SETUP:
#
#   systemctl --user start metapi.service
#   systemctl --user enable metapi.service   # (optional: auto-start at login)
#
#   Then open http://127.0.0.1:4000, log in with AUTH_TOKEN, and add your upstream
#   sites/accounts and downstream keys through the UI. Config lives in the
#   metapi-data volume (SQLite) — no external database required.
#
# CLIENT SETUP (after creating a downstream key in the UI):
#
#   claude-code  → ANTHROPIC_BASE_URL=http://127.0.0.1:4000
#                  ANTHROPIC_AUTH_TOKEN=<your PROXY_TOKEN or a downstream key>
#   codex-cli    → base_url = "http://127.0.0.1:4000/v1" in ~/.codex/config.toml

[Unit]
Description=metapi unified LLM gateway
Documentation=https://github.com/cita-777/metapi
After=network-online.target
Wants=network-online.target

[Container]
# Pinned by digest (image built 2026-06-24, MIT-licensed, upstream: github.com/cita-777/metapi).
# The image lives under a personal Docker Hub account, so pin rather than track :latest.
# TODO(kai): self-build metapi from source and repoint this at your own registry — see
# docs/self-build-todo.md. To bump: skopeo inspect docker://docker.io/1467078763/metapi:latest
Image=docker.io/1467078763/metapi@sha256:d6118229e7d2423262b253a419baf18c22f1682a7bbc7d3c756d090aa2b295c6
ContainerName=metapi
PublishPort=127.0.0.1:4000:4000
# :U chowns the volume to the container's run-user (image runs as uid 1000
# "node", not root). Without it the node user can't create the SQLite DB in
# a root-owned named volume → better-sqlite3 SQLITE_CANTOPEN at startup.
Volume=metapi-data:/app/data:U
Environment=PORT=4000
Environment=DATA_DIR=/app/data
Environment=TZ=Asia/Hong_Kong
# Quote cron values: systemd's Environment= parser splits on unquoted
# whitespace, so `CHECKIN_CRON=0 8 * * *` would reach the container as
# `CHECKIN_CRON=0` and node-cron crashes on the one-field expression.
Environment="CHECKIN_CRON=0 8 * * *"
Environment="BALANCE_REFRESH_CRON=0 * * * *"

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
